# deployment/jper.conf
[program:jper]
command=/home/green/jper/bin/gunicorn -c /home/green/jper/src/jper/deployment/gconf.py service.web:app
user=green
directory=/home/green/jper/src/jper
stdout_logfile=/var/log/supervisor/%(program_name)s-access.log
stderr_logfile=/var/log/supervisor/%(program_name)s-error.log
autostart=true
autorestart=true
stopasgroup=true
environment=PYTHONPATH=/home/green/jper/src/jper/OAUtils/src/
